# myfitnesspal-restapi.service
[Unit]
Description=MyFitnesspal API

[Service]
Type=simple
ExecStart=/bin/bash /opt/myfitnesspal-api/run.sh

[Install]
WantedBy=multi-user.target
